# archlive/airootfs/etc/systemd/system/ollama.service.d/override.conf
[Unit]
Wants=ai-os-ollama-storage.service
After=ai-os-ollama-storage.service

[Service]
Environment=OLLAMA_MODELS=/var/lib/ollama/models
